# deploy/systemd/jarvis-mod-property-intel.service
[Unit]
Description=JARVIS v2 property-intel module
Documentation=https://github.com/swissdog/property-intel
After=network-online.target jarvis-core-api.service
Wants=network-online.target

[Service]
Type=simple
User=sami
Group=sami
WorkingDirectory=/home/sami/property-intel
EnvironmentFile=-/home/sami/property-intel/.env
ExecStart=/home/sami/.local/bin/uv run --no-sync uvicorn jarvis_property_intel.main:app --host 0.0.0.0 --port 8031

Restart=on-failure
RestartSec=5
TimeoutStartSec=120

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
